# GToM/deploy/systemd/gtom.service
[Unit]
Description=GToM - Conflict prediction service for G-Stack
After=network.target postgresql.service

[Service]
Type=simple
User=gtom
Group=gtom
WorkingDirectory=/opt/gtom
Environment=NODE_ENV=production
Environment=PORT=3003
Environment=HEALTH_PORT=8080
EnvironmentFile=/etc/gtom/gtom.env
ExecStart=/usr/bin/node /opt/gtom/dist/serve.js
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30
KillSignal=SIGTERM
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gtom

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/gtom /var/log/gtom

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
